# .config/systemd/system/user@.service.d/local.conf
[Service]
Environment="PATH=/home/juuso/bin:/usr/lib/bin:/usr/local/bin:/usr/bin:/bin"
Environment="EDITOR=nvim"
Environment="BROWSER=brave"
